# deb/supervisor-config/supervisor-sync.conf
[program:sync_content]
command=/scripts/sync_content.sh
numprocs=1
autostart=true
autorestart=true
